# cldperf-nflx-lab-benchmarks-main/systemctl-service/run-benchmark.service
[Unit]
Description=Autobench - Run Benchmarks
After=nflx.target

[Service]
Type=oneshot
ExecStart=/usr/share/autobench/run-benchmarks

[Install]
WantedBy=multi-user.target
